# rootfs/etc/s6-overlay/s6-rc.d/internxt-webdav/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

# Get configuration values
EMAIL=$(bashio::config 'email')
PASSWORD=$(bashio::config 'password')
WEBDAV_PORT=$(bashio::config 'webdav_port')
WEBDAV_PROTOCOL=$(bashio::config 'webdav_protocol')

bashio::log.info "Logging into Internxt..."

# Build login command
LOGIN_CMD="internxt login --email '${EMAIL}' --password '${PASSWORD}'"

# Add 2FA if configured
if bashio::config.has_value 'two_factor_code' && [ -n "$(bashio::config 'two_factor_code')" ]; then
    TWO_FACTOR=$(bashio::config 'two_factor_code')
    LOGIN_CMD="${LOGIN_CMD} --twofactor '${TWO_FACTOR}'"
fi

# Execute login
eval "${LOGIN_CMD}"

if [ $? -ne 0 ]; then
    bashio::log.error "Failed to login to Internxt"
    exit 1
fi

bashio::log.info "Starting Internxt WebDAV server on port ${WEBDAV_PORT} (${WEBDAV_PROTOCOL})"

# Start WebDAV server
exec internxt webdav --port "${WEBDAV_PORT}" --protocol "${WEBDAV_PROTOCOL}"
